# deploy/server/health-app.service
[Unit]
Description=Health Center Backend Service
Documentation=https://github.com/health-center
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/health-center
ExecStart=/usr/bin/java -jar /opt/health-center/backend.jar \
    --spring.profiles.active=prod \
    --server.port=8080
ExecStop=/bin/kill -15 $MAINPID
Restart=always
RestartSec=10
StandardOutput=append:/opt/health-center/logs/console.log
StandardError=append:/opt/health-center/logs/error.log

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/health-center/uploads
ReadWritePaths=/opt/health-center/logs

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=1G

[Install]
WantedBy=multi-user.target
